# scripts/trusttunnel.service.template
# Service template for the TrustTunnel endpoint
# Adjust paths at `ExecStart` and copy file to systemd system directory:
#
# cp trusttunnel.service.template /etc/systemd/system/trusttunnel.service
# systemctl daemon-reload
# systemctl enable --now trusttunnel

[Unit]
Description=TrustTunnel endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trusttunnel
ExecStart=/opt/trusttunnel/trusttunnel_endpoint vpn.toml hosts.toml
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
